# ansible/roles/services/files/aoe.ini
[program:aoe]
command=/usr/bin/npm start
directory=/opt/sources
autostart=true
autorestart=true
environment=NODE_ENV=dev
stdout_logfile=/var/log/supervisor/aoe.log
stderr_logfile=/var/log/supervisor/aoe.err
user=vagrant
